# supervisord/supervisor-rsyslog.conf
[program:rsyslog]
#command=/bin/bash -c ". /etc/default/rsyslog && rsyslogd -n $RSYSLOGD_OPTIONS"
command=rsyslogd -n -c5
