# widgets/supervisor-widgets.conf
[program:widgets]
command=/widget.sh %(process_num)02d
numprocs=3
autostart=true
autorestart=true
redirect_stderr=true
process_name=%(program_name)s-%(process_num)s 
